# ect/upstart/airbnb.conf
# upstart script for /etc/init/airbnb.conf
# sudo  start airbnb
# debug: $ init-checkconf config/airbnb.conf

description "gunicorn application server running AirBnB_clone"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data

